# misc/cbc_lifecycle/cbc_lifecycle.service
[Unit]
Description=CBC lifecycle service

[Service]
Type=simple
ExecStart=/usr/bin/cbc_lifecycle
Restart=no

[Install]
WantedBy=multi-user.target
