# misc/cbc_lifecycle/cbc_lifecycle.service
[Unit]
Description=CBC lifecycle service

[Service]
Type=simple
ExecStart=/usr/bin/cbc_lifecycle
Restart=no
ExecStop=/usr/bin/killall -s TERM cbc_lifecycle

[Install]
WantedBy=multi-user.target
