# misc/debug_tools/acrn_log/acrnlog.service
[Unit]
Description=ACRN hypervisor log

[Service]
Type=simple
ExecStart=/usr/bin/acrnlog

[Install]
WantedBy=multi-user.target

