# misc/services/life_mngr/life_mngr.service.in
[Unit]
Description=ACRN lifemngr daemon
After=setserial.service
After=network.target
After=systemd-logind.service

[Service]
Type=simple
ExecStart=@bindir@/life_mngr
StandardOutput=journal
StandardError=journal
Restart=on-failure

[Install]
WantedBy=multi-user.target
